# templates/crontab
*/10 * * * * %(dir)s/start_supervisor.sh
