# templates/gunicorn.conf
[program:%(project)s]
command=%(virtualenv)s/bin/gunicorn_django -w 2 -b 0.0.0.0:%(port)s
directory=%(project_dir)s
user=%(user)s
autostart=true
autorestart=true
redirect_stderr=True
